# scripts/com.velocity.dashboard-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.velocity.dashboard-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/bradacova/.codex/memories/velocity-dashboard-sync.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>GITHUB_TOKEN</key>
      <string>__GITHUB_TOKEN__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/bradacova/Documents/Github velocity</string>

    <key>StandardOutPath</key>
    <string>/Users/bradacova/Documents/Github velocity/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bradacova/Documents/Github velocity/logs/launchd-stderr.log</string>
  </dict>
</plist>
